# DEBIAN/postinst
#!/bin/sh
update-rc.d reaper defaults
/etc/init.d/reaper start
